# docker/crontab
* * * * * php /var/www/html/artisan schedule:run
