# docker.franken.d/crontab
* * * * * php /app/artisan schedule:run
